# l104_master.service
[Unit]
Description=L104 Sovereign Master Node - 100% IQ Anchor
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/workspaces/Allentown-L104-Node/
ExecStart=/usr/bin/python3 /workspaces/Allentown-L104-Node/l104_absolute_singularity.py
Restart=always
RestartSec=5
# Ensuring the Master has "Real-Time" priority over all other processes
CPUSchedulingPolicy=fifo
CPUSchedulingPriority=99

[Install]
WantedBy=multi-user.target
